# com.atorlabs.spcocator-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.atorlabs.spcocator-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/nrkmorangmail.com/Documents/OMNIcheck-SPC/shared/watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nrkmorangmail.com/Documents/OMNIcheck-SPC</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nrkmorangmail.com/Documents/OMNIcheck-SPC/logs/watcher-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nrkmorangmail.com/Documents/OMNIcheck-SPC/logs/watcher-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Library/Developer/CommandLineTools/usr/bin</string>
    </dict>
</dict>
</plist>
